# MyModule/UniCron/UniCron.cron
* * * * * /data/adb/modules/UniCron/test.sh